# root/etc/services.d/facette/run
#!/usr/bin/with-contenv bash

USER_NAME=facette

exec 2>&1 \
        s6-setuidgid ${USER_NAME} /usr/local/bin/facette
